# files/24.10/etc/uci-defaults/360T7
# 24.10 openwrt uci一键配置

# /etc/config/autoreboot
uci set autoreboot.cfg016bf2.enabled='1'
uci set autoreboot.cfg016bf2.week='3'
uci set autoreboot.cfg016bf2.hour='9'
uci commit autoreboot

# /etc/config/firewall
uci del firewall.cfg02dc81.network
uci add_list firewall.cfg02dc81.network='lan'
uci del firewall.cfg03dc81.network
uci add_list firewall.cfg03dc81.network='wan'
uci add_list firewall.cfg03dc81.network='wan6'
uci add_list firewall.cfg03dc81.network='wwan'
uci commit firewall

# /etc/config/network
uci set network.wwan=interface
uci set network.wwan.proto='dhcp'
uci commit network

# /etc/config/wireless
uci set wireless.wifinet2=wifi-iface
uci set wireless.wifinet2.device='MT7981_1_1'
uci set wireless.wifinet2.mode='sta'
uci set wireless.wifinet2.network='wwan'
uci set wireless.wifinet2.ssid='Hjx'
uci set wireless.wifinet2.encryption='psk2'
uci set wireless.wifinet2.key='HjxWpy2580'
uci del wireless.MT7981_1_2.noscan
#uci set wireless.default_MT7981_1_2.hidden='1'
uci set wireless.default_MT7981_1_2.encryption='sae-mixed'
uci set wireless.default_MT7981_1_2.key='88888888..'
#uci set wireless.default_MT7981_1_1.hidden='1'
uci set wireless.default_MT7981_1_1.encryption='sae-mixed'
uci set wireless.default_MT7981_1_1.key='88888888..'
uci set wireless.default_MT7981_1_1.disabled='1'
uci commit wireless

# /etc/config/system
uci set system.led_network.dev='apcli0'
uci del system.led_network.mode
uci add_list system.led_network.mode='link'
uci commit system